# scheduler/bus-gps-upd.timer
[Unit]
Description=Run Bus GPS API Loader Every 10 Seconds (Calendar-based)
Requires=bus-gps-upd.service

[Timer]
# Run at 0, 10, 20, 30, 40, 50 seconds of each minute
OnCalendar=*:*:0,10,20,30,40,50

# Make timing as accurate as possible
AccuracySec=1us

# Run immediately if we missed a scheduled time
Persistent=true

[Install]
WantedBy=timers.target
