# scheduler/bus-gps-upd.timer
[Unit]
Description=Run Bus GPS API Loader Every 10 Seconds (Calendar-based)
Requires=bus-gps-upd.service

[Timer]
# Run at 0, 10, 20, 30, 40, 50 seconds of each minute
# Between 22:00 UTC and 19:59 UTC next day (3AM-1AM GMT+5)
OnCalendar=*-*-* 22,23,00,01,02,03,04,05,06,07,08,09,10,11,12,13,14,15,16,17,18,19:*:0,10,20,30,40,50

# Make timing as accurate as possible
AccuracySec=1us

# Run immediately if we missed a scheduled time
Persistent=true

[Install]
WantedBy=timers.target